# test/app/service/sshd/run
#!/bin/sh

/usr/bin/ssh-keygen -A

exec /usr/sbin/sshd -D -f ./sshd_config
